# fabric_bolt/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:nginx]
command=/usr/sbin/nginx

[program:fabric_bolt]
command=fabric-bolt runserver 5000


